# fileserver.conf
[program:device-controller]
user=root
command=/home/pi/idol/fileserver
autostart = true
autorestart = true

